# cfengine/chef/templates/default/supervisord.conf.erb
[supervisord]
http_port=/var/tmp/supervisor.sock
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=debug
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[supervisorctl]
serverurl=unix:///var/tmp/supervisor.sock

[program:phantomjs]
directory = /tmp
user = root
command = /opt/phantomjs-1.9.8-linux-x86_64/bin/phantomjs --webdriver=8910
;redirect_stderr=true
logfile=/var/log/supervisor/phantomjs.log
autostart=true
autorestart=true